# surfnet/root/etc/services.d/xrayr/run
#!/usr/bin/with-contenv bash

exec \
"XrayR", "--config", "/etc/XrayR/config.yml"
